# dotfiles/system_Linux/config/systemd/user/mako.service.d/override.conf
[Unit]
PartOf=
BindsTo=sway-session.target
After=
ConditionEnvironment=WAYLAND_DISPLAY

[Install]
WantedBy=
WantedBy=sway-session.target
